# com.example.xhs-topics-heat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ocean.xhs-topics-heat</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOUR_USER/xhs-topics/venv/bin/python</string>
        <string>/Users/YOUR_USER/xhs-topics/heat_refresh.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>9</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/YOUR_USER/Library/Logs/xhs-topics-heat.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USER/Library/Logs/xhs-topics-heat.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HTTP_PROXY</key><string>http://127.0.0.1:7897</string>
        <key>HTTPS_PROXY</key><string>http://127.0.0.1:7897</string>
        <key>LANG</key><string>en_US.UTF-8</string>
        <key>LC_ALL</key><string>en_US.UTF-8</string>
    </dict>
</dict>
</plist>
